# systemd/dictate.service
[Unit]
Description=voxtral-dictate daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/hobby/voxtral-dictate/dictate daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
